# doca-dpu-repo-ubuntu2204-local/mlnx-dpdk-dev_22.11.0-2401.1.0.24010330.2401.1.0_arm64.deb_dir/etc/profile.d/mlnx-dpdk-aarch64-linux-gnu.sh
export PKG_CONFIG_PATH=${PKG_CONFIG_PATH}:/opt/mellanox/dpdk/lib/aarch64-linux-gnu/pkgconfig
