# doca-dpu-repo-ubuntu2204-local/vrf_0.9.12-1_arm64.deb_dir/etc/profile.d/vrf.sh
VRF=$(ip vrf identify)
[ -n "${VRF}" ] && VRF=":${VRF}"
export VRF

PS1='\u@\h${VRF}:\w\$ '
