# doca-dpu-repo-ubuntu2204-local/libgrpc-dev_1.39.0-1_arm64.deb_dir/etc/profile.d/libgrpc.sh
export PATH=${PATH}:/opt/mellanox/grpc/bin
export PKG_CONFIG_PATH=${PKG_CONFIG_PATH}:/opt/mellanox/grpc/lib/pkgconfig
